# diamond/supervisord.conf
[supervisord]
nodaemon=true

[program:diamond]
command=/etc/init.d/diamond start
